# streamer/tijdvorm-receiver.service
[Unit]
Description=Tijdvorm Stream Receiver
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/home/sjoerd/tijdvorm/streamer
ExecStart=/usr/bin/xinit /home/sjoerd/tijdvorm/streamer/start_receiver.sh -- :0 -nocursor
Restart=always
RestartSec=5

# Environment defaults (override in /etc/default/tijdvorm-receiver)
Environment=STREAM_URL=http://m4.netlob:8000/stream
EnvironmentFile=-/etc/default/tijdvorm-receiver

[Install]
WantedBy=multi-user.target
